# templates/serverkit-backend.service.in
[Unit]
Description=ServerKit Backend API
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=@USER@
Group=@USER@
WorkingDirectory=@SERVERKIT_DIR@/backend
Environment="PATH=@SERVERKIT_VENV_DIR@/bin:/usr/local/bin:/usr/bin:/bin"
Environment="FLASK_ENV=production"
EnvironmentFile=@SERVERKIT_DIR@/.env
ExecStart=@SERVERKIT_VENV_DIR@/bin/gunicorn -k geventwebsocket.gunicorn.workers.GeventWebSocketWorker -w 1 -b 0.0.0.0:@PORT@ --timeout 120 --access-logfile @LOG_DIR@/access.log --error-logfile @LOG_DIR@/error.log run:app
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5

# Security settings (while still allowing system management)
NoNewPrivileges=no
ProtectSystem=no
ProtectHome=no

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=serverkit

[Install]
WantedBy=multi-user.target
